# demon/test_demon.service.example
[Unit]
Description=Test daemon

[Service]
User=<your_user_name>
ExecStart=sh <path_to_project>/test_task_carbonSoft/demon.sh


